# deploy/ochenstarik-smm-provisioning-helper.service
[Unit]
Description=Ochenstarik Server Monitor Manager restricted provisioning helper
After=local-fs.target

[Service]
Type=simple
User=root
Group=ochenstarik-smm-agent
UMask=0007
RuntimeDirectory=ochenstarik-server-monitor-manager
RuntimeDirectoryMode=0750
WorkingDirectory=/usr/local/lib/ochenstarik-server-monitor-manager/provisioning-helper
ExecStart=/usr/local/lib/ochenstarik-server-monitor-manager/provisioning-helper/ochenstarik-smm-provisioning-helper
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
